# cm5_autostart/ld19-stack.service
[Unit]
Description=LD19 LiDAR stack (preprocess + monitor, driver provided by ld19.service)
After=ld19.service
Requires=ld19.service

[Service]
Type=simple
User=tourrobotsub
Environment=ROS_DOMAIN_ID=0
Environment=ROS_LOCALHOST_ONLY=0
ExecStart=/usr/local/bin/start_ld19_stack.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
